# com.therepute.googlenews.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.therepute.googlenews.monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/willvalentine/Soup_Google_News_Ingestion/-Users-willvalentine-the_soup-google_news_ingestion_layer/monitor_service.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/willvalentine/Soup_Google_News_Ingestion/-Users-willvalentine-the_soup-google_news_ingestion_layer</string>
    
    <key>StandardOutPath</key>
    <string>/Users/willvalentine/Soup_Google_News_Ingestion/-Users-willvalentine-the_soup-google_news_ingestion_layer/logs/monitor.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/willvalentine/Soup_Google_News_Ingestion/-Users-willvalentine-the_soup-google_news_ingestion_layer/logs/monitor_error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StartInterval</key>
    <integer>300</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist> 